# rootfs/etc/s6-overlay/s6-rc.d/webssh2-log/run
#!/command/execlineb -P
s6-setuidgid webssh2
s6-log -b n20 s1000000 /var/log/webssh2
