# rootfs/etc/s6-overlay/s6-rc.d/webssh2/run
#!/bin/bash

# WebSSH2 service script for s6-overlay
cd /usr/src/webssh2

# Source environment variables created by initialization script
if [[ -f /etc/webssh2-env ]]; then
    source /etc/webssh2-env
fi

# Run as webssh2 user - s6-setuidgid preserves environment variables
exec s6-setuidgid webssh2 node index.js