# rootfs/etc/s6-overlay/s6-rc.d/webssh2/run
#!/command/execlineb -P
foreground { /usr/local/bin/webssh2-setup.sh }
s6-envdir /etc/s6-overlay/s6-rc.d/webssh2/env
cd /usr/src/webssh2
s6-setuidgid webssh2
fdmove -c 2 1
/usr/bin/node dist/index.js
