# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/bin/bash

# NGINX service script for s6-overlay
exec nginx -g 'daemon off;'